# build/phantom-overlay/airootfs/etc/systemd/system/phantom-gpu-setup.service
[Unit]
Description=Phantom OS GPU Auto-Detection and Setup
After=local-fs.target
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/phantom-gpu-setup
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
